# ops/dms-tunnel.service
# Persistent SSH tunnel from worksserver → dmsserver MySQL.
# Install on worksserver as /etc/systemd/system/dms-tunnel.service.
#
# The LMIS backend container reaches MySQL at host.docker.internal:3307
# (the loopback side of this tunnel). dmsserver only exposes 3306 on
# 127.0.0.1, which is why the tunnel is needed.
#
# Install:
#   sudo install -m 0644 ops/dms-tunnel.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now dms-tunnel.service
#   sudo systemctl status dms-tunnel.service
#
# Verify from worksserver:
#   ss -tln | grep 3307
#   mysql -h 127.0.0.1 -P 3307 -udmsuser -p dmsworkflow -e 'SELECT 1'
#
# Verify from backend container:
#   docker compose -f /opt/bda-lmis/docker-compose.prod.yml exec backend \
#       python -c "import socket; s=socket.socket(); s.connect(('host.docker.internal',3307)); print('OK')"
#
# (host.docker.internal resolves on Linux via the extra_hosts entry on
# the backend service in docker-compose.prod.yml — see the deploy doc.)

[Unit]
Description=SSH tunnel to dmsserver MySQL (127.0.0.1:3307 -> dms:3306)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=itadmin
ExecStart=/usr/bin/ssh -N \
    -o ServerAliveInterval=30 \
    -o ServerAliveCountMax=3 \
    -o ExitOnForwardFailure=yes \
    -o StrictHostKeyChecking=accept-new \
    -i /home/itadmin/.ssh/dms-tunnel \
    -L 127.0.0.1:3307:127.0.0.1:3306 \
    -p 2222 \
    itadmin@69.62.80.7
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
