# ops/dms-tunnel.service
# Persistent SSH tunnel from worksserver → dmsserver MySQL.
# Install on worksserver as /etc/systemd/system/dms-tunnel.service.
#
# Bind address note: the tunnel listens on 172.17.0.1:3307 (the docker0
# bridge gateway), NOT 127.0.0.1:3307. Reason: the LMIS backend lives in
# a docker bridge network and reaches the host via 172.17.0.1 (mapped
# to `host.docker.internal` via the extra_hosts entry in
# docker-compose.prod.yml). Binding to 127.0.0.1 makes the tunnel only
# reachable from the host's own loopback — useless to the container.
# Binding to 172.17.0.1 keeps it off the public network while letting
# the docker bridge reach it.
#
# Install:
#   sudo install -m 0644 ops/dms-tunnel.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now dms-tunnel.service
#   sudo systemctl status dms-tunnel.service
#
# Verify on worksserver:
#   ss -tln | grep 3307        # expect LISTEN on 172.17.0.1:3307
#
# Verify from backend container:
#   docker compose -f /opt/bda-lmis/docker-compose.prod.yml exec backend \
#       python -c "import socket; s=socket.socket(); s.settimeout(3); \
#                  s.connect(('host.docker.internal',3307)); print('OK')"

[Unit]
Description=SSH tunnel to dmsserver MySQL (172.17.0.1:3307 -> dms:3306)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=itadmin
ExecStart=/usr/bin/ssh -N \
    -o ServerAliveInterval=30 \
    -o ServerAliveCountMax=3 \
    -o ExitOnForwardFailure=yes \
    -o StrictHostKeyChecking=accept-new \
    -i /home/itadmin/.ssh/dms-tunnel \
    -L 172.17.0.1:3307:127.0.0.1:3306 \
    -p 2222 \
    itadmin@69.62.80.7
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
